# supervisord.conf
[program:node]
command = node app.js
directory = /home/dotcloud/current
NODE_ENV = production